# dataspace-startup.service
[Unit]
Description=Run DataSpace startup script at boot
After=network.target

[Service]
Type=oneshot
User=dataspace
WorkingDirectory=/opt/dataspace/Platform
ExecStart=/opt/dataspace/Platform/startup.sh
ExecStop=/usr/bin/docker compose -f /opt/dataspace/Platform/docker-compose.yml down

[Install]
WantedBy=multi-user.target
